# secure-ai-chat.service
[Unit]
Description=Secure AI Chat Application
After=network.target

[Service]
Type=simple
User=adminuser
Group=adminuser
WorkingDirectory=/home/adminuser/secure-ai-chat
Environment="NODE_ENV=production"
Environment="PORT=3000"
Environment="HOSTNAME=0.0.0.0"

# Ensure nvm-installed Node v25.2.1 is used
# IMPORTANT: Update the path below to match your actual nvm Node installation path
# Find it with: which node (after running: nvm use 25.2.1)
# Example: /home/adminuser/.nvm/versions/node/v25.2.1/bin/npm
# Note: If Node path changes, update this line accordingly
ExecStart=/home/adminuser/.nvm/versions/node/v25.2.1/bin/npm start
# Alternative: Use absolute node path: ExecStart=/home/adminuser/.nvm/versions/node/v25.2.1/bin/node /home/adminuser/secure-ai-chat/node_modules/.bin/next start
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=secure-ai-chat

# Security
NoNewPrivileges=true
PrivateTmp=true
ReadWritePaths=/home/adminuser/secure-ai-chat/.secure-storage /home/adminuser/secure-ai-chat/.next

# Health check via systemd timer (optional - can also use external monitoring)
# Health check endpoint: http://localhost:3000/api/health

[Install]
WantedBy=multi-user.target
