# scripts/deploy/secure-ai-chat.service
[Unit]
Description=Secure AI Chat Application
After=network.target

[Service]
Type=simple
User={{APP_USER}}
Group={{APP_USER}}
WorkingDirectory={{APP_DIR}}
EnvironmentFile=/etc/secure-ai-chat.env

# Ensure NODE_ENV=production (override env file if needed)
Environment="NODE_ENV=production"

# Load nvm and use correct Node.js version, then validate startup and start
# Validate startup ensures directories exist with correct permissions
ExecStartPre=/usr/bin/bash -c 'cd {{APP_DIR}} && bash scripts/validate-startup.sh || exit 1'
ExecStart=/usr/bin/env bash -lc 'source "$HOME/.nvm/nvm.sh" && nvm use {{NODE_VERSION}} && cd {{APP_DIR}} && npm start'

Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal
SyslogIdentifier=secure-ai-chat

# Security
NoNewPrivileges=true
PrivateTmp=true
ReadWritePaths={{APP_DIR}}/.secure-storage {{APP_DIR}}/.next {{APP_DIR}}/.storage

[Install]
WantedBy=multi-user.target
